# oscam/install/files/etc/service/oscam/run
#!/bin/bash

chown -R "$PUSER":"$PGROUP" /config /log
exec /sbin/setuser "$PUSER" /usr/bin/oscam

